# test_monitoring.timer
[Unit]
Description=timer for start proccess

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
