# test_monitoring.service
[Unit]
Description=test_monitoring

[Service]
Type=simple
ExecStart=/usr/local/bin/test_monitoring.sh
Restart=no
